# config/com.ebay.clearserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ebay.clearserver</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/samwhite/Documents/ebay-search/venv/bin/python</string>
        <string>/Users/samwhite/Documents/ebay-search/clear_server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/samwhite/Documents/ebay-search</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/samwhite/Documents/ebay-search/clear_server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/samwhite/Documents/ebay-search/clear_server.log</string>
</dict>
</plist>
